# packages/palhm/palhm-0.0.3.tar.gz/palhm-0.0.3/src/palhm/conf/crontab
# PALHM will produce stderr on exception. Mail the output to root
MAILTO="root"
# Run default task every Sunday at midnight
0  0  *  *  sun root python -m palhm -q run
# Check dnssec validity every hour
# 0  *  *  *  *   root systemd-run -qP -p User=palhm -p Nice=15 -p ProtectSystem=strict -p ReadOnlyPaths=/ -p PrivateDevices=true --wait python -m palhm -q run check-dnssec
