# packages/palhm/palhm-0.0.3.tar.gz/palhm-0.0.3/src/palhm/conf/palhm-boot-report.service
[Unit]
Description=PALHM send boot report mail
After=postfix.service sendmail.service exim.service dovecot.service network-online.target

[Service]
Type=oneshot
ExecStart=/usr/bin/python -m palhm -q boot-report
Nice=10

[Install]
WantedBy=multi-user.target
